# cmos-reset/xyz.openbmc_project.Software.Host.Updater0.service.in
[Unit]
Description=platform cmos reset service
Before=xyz.openbmc_project.FruDevice.service

[Service]
ExecStart=/usr/libexec/platform-init/cmos-reset @PLATFORM_NAME@
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Software.Host.Updater0

[Install]
WantedBy=multi-user.target
